# launchd/com.turublok.lock.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.turublok.lock</string>

    <key>ProgramArguments</key>
    <array>
        <string>__BINARY_PATH__</string>
        <string>--lock</string>
    </array>

    <!-- Trigger di setiap jam dalam lock window (23:00 - 06:00).
         Kalau Mac di-force-shutdown di tengah malem, jam berikutnya launchd nyalain lagi. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>1</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <!-- Trigger pas login (termasuk setelah restart). Binary self-check window + sentinel. -->
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>__HOME__/Projects/TURU-BLOK/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Projects/TURU-BLOK/logs/launchd.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>Asia/Jakarta</string>
    </dict>
</dict>
</plist>
