# launchd/com.turublok.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.turublok.web</string>

    <!-- Control server (PWA buat HP istri). Selalu nyala selama Mac hidup. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>__WEB_SERVER__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/Projects/TURU-BLOK/logs/web.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Projects/TURU-BLOK/logs/web.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>Asia/Jakarta</string>
    </dict>
</dict>
</plist>
